# qr3_96.conf
[program:qr3_96]
command=python3 /home/curator3/qr3_96/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_96/test.err.log
stdout_logfile=/home/curator3/qr3_96/test.err.log
